# deploy/systemd/emigro-portugal-community.timer
[Unit]
Description=Emigro Portugal community daily at 07:00 UTC

[Timer]
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
